# build/meta-rootfs/recipes-rootfs/linux/files/board/common/rootfs_overlay/etc/systemd/system/fb_mapper.service
[Unit]
Description=Starting Framebuffer mapper service

[Service]
Restart=always
ExecStart=/root/fb_mapper

[Install]
WantedBy=multi-user.target
