# build/meta-rootfs/recipes-rootfs/linux/files/board/common/rootfs_overlay/etc/systemd/system/input_forwarding.service
[Unit]
Description=Starting input forwarding service

[Service]
Restart=always
ExecStart=/root/input_forwarding

[Install]
WantedBy=multi-user.target
