# build/meta-rootfs/recipes-rootfs/linux/files/board/common/rootfs_overlay/etc/systemd/system/soo.service
[Unit]
Description=Smart Object Oriented

[Service]
Type=o eshot
ExecStart=/bin/sh -c "mkdir /dev/soo && mknod /dev/soo/core c 126 0"
TimeoutStartSec=0
RemainAfterExit=yes
StandardError=journal
StandardOutput=journal
StandardInput=null


[Install]
WantedBy=default.target
