# build/meta-rootfs/recipes-rootfs/linux/files/board/common/rootfs_overlay/etc/systemd/system/test.service
[Unit]
Description=Test Service

[Service]
Type=simple
ExecStart=/root/test_service
TimeoutStartSec=0
RemainAfterExit=yes
StandardError=journal
StandardOutput=journal
StandardInput=null


[Install]
WantedBy=default.target